# crontab
0 */2 * * * root /backup.sh > /var/log/cron.log 2>&1
